# conf/etc/periodic/15min/01-vtiger
#!/bin/sh

run(){
    eval "su - nginx -s /bin/ash -c 'cd /var/www/drupal/ && "$@"'" 
    return $? 
}
run php /var/www/vtigercrm/vtigercron.php
